# supervisord.conf
[supervisord]
nodaemon=true
#loglevel=warn

[program:serposcope]
command=/bin/bash /docker-entrypoint.sh
startsecs = 1
autostart = true